# docker/etc/services.d/scal/run
#!/usr/bin/execlineb -P

/usr/share/scal/bin/scal-server -p 80 -s /var/lib/dlubm -t /tmp &> /var/lib/dlubm/scal
